# roles/backup/templates/vaultwarden-backup.service.j2
[Unit]
Description=Vaultwarden backup
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-vaultwarden
